# provisioning/virtual_machine/rootfs/etc/systemd/system/runatstartup.service
[Unit]
Description=run script(s) at startup

[Service]
WorkingDirectory=/
ExecStart=/usr/bin/runatstartup
Restart=no

[Install]
WantedBy=multi-user.target
