# provisioning/linux_containers/common/rootfs/etc/systemd/system/runatstartup.service
[Unit]
Description=run script(s) at startup

[Service]
WorkingDirectory=/
ExecStart=/sbin/runatstartup
Restart=no

[Install]
WantedBy=multi-user.target
